# sway-atomic-files/wallpaper-cycle.timer
[Unit]
Description=Change the wallpaper every 4 hours

[Timer]
OnBootSec=4h
OnUnitActiveSec=4h

[Install]
WantedBy=timers.target
